# deploy/grid-monitor.service.in
[Unit]
Description=Electricity Grid Monitor
After=local-fs.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__PROJECT_DIR__
EnvironmentFile=__PROJECT_DIR__/.env
ExecStart=/usr/bin/python3 -m grid_monitor monitor
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
